# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-sdkman-add-package/run
#!/usr/bin/with-contenv bash

echo "\
    curl \
    zip \
    unzip \
    sed" >> /mod-repo-packages-to-install.list
